# supervisord.conf
[supervisord]
nodaemon=true

[program:holed]
command=holed -addr=tcp://0.0.0.0:5188
