# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:holed]
command=holed -addr=tcp://0.0.0.0:5188

[program:ssserver]
command=ssserver -p 5189 -k 6a3503861efaf6c -d start